# target/linux/generic/base-files/etc/uci-defaults/99_custom
uci set network.lan.ipaddr=192.168.2.1
uci commit network

uci set system.@system[0].timezone=CST-8
uci set system.@system[0].zonename=Asia/Shanghai
uci commit system

uci set luci.main.lang=zh_cn
uci commit luci

echo net.ipv6.conf.all.disable_ipv6=1 >> /etc/sysctl.conf
sysctl -p

exit 0
